# roles/overlay_network/files/overlay-network.service
[Unit]
Description=Overlay network service
After=network-online.target wg-quick.target
Wants=network-online.target wg-quick.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/overlay-network/overlay-network start
ExecStop=/opt/overlay-network/overlay-network stop

[Install]
WantedBy=multi-user.target
